# centreon-dsm/packaging/debian/centreon-dsm-server.postinst
#!/bin/sh

if [ "$1" = "configure" ] ; then

  systemctl daemon-reload
  systemctl enable dsmd.service

fi
exit 0