# centreon/packaging/debian/centreon-common.postinst
#!/bin/sh

if [ "$1" = "configure" ] ; then
  chown -R centreon:centreon \
    /etc/centreon \
    /var/cache/centreon
  chmod -R 0775 /var/cache/centreon
  chmod -R 2775 /var/cache/centreon/config
fi

# Initial installation
systemctl --no-reload preset centreon.service || : &>/dev/null || :
